# simulation-daemon.service
[Unit]
Description=AI Foundry Agent Simulation Daemon
After=network.target

[Service]
Type=simple
User=REPLACE_ME
Group=REPLACE_ME
WorkingDirectory=/path/to/this/repo
ExecStart=/usr/bin/env python3 /path/to/this/repo/simulation_daemon.py --config /path/to/this/repo/simulation_daemon_config.json

# Restart policy
Restart=on-failure
RestartSec=30

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=append:/path/to/this/repo/daemon_results/daemon_stdout.log
StandardError=append:/path/to/this/repo/daemon_results/daemon_stderr.log

[Install]
WantedBy=multi-user.target
