# simulation-daemon.service
[Unit]
Description=AI Foundry Agent Simulation Daemon
After=network.target

[Service]
Type=simple
User=azureuser
Group=azureuser
WorkingDirectory=/home/azureuser/code/AI-Foundry-Agent-Simulation
Environment="PATH=/home/azureuser/code/AI-Foundry-Agent-Simulation/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/azureuser/code/AI-Foundry-Agent-Simulation/.venv/bin/python3 /home/azureuser/code/AI-Foundry-Agent-Simulation/simulation_daemon.py

# Restart policy
Restart=on-failure
RestartSec=30

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=append:/home/azureuser/code/AI-Foundry-Agent-Simulation/daemon_results/daemon_stdout.log
StandardError=append:/home/azureuser/code/AI-Foundry-Agent-Simulation/daemon_results/daemon_stderr.log

[Install]
WantedBy=multi-user.target
